# sojson.ini
[program:sojson]
command=/Users/mico/Github/sojson/sojson
directory=/Users/mico/Github/sojson
autostart=true
autorestart=true
stderr_logfile=/Users/mico/Github/sojson/log/sojson.err.log
stdout_logfile=/Users/mico/Github/sojson/log/sojson.out.log
user=mico
environment=PATH="/usr/local/bin:/usr/bin:/bin"
numprocs=1
startsecs=3
stopwaitsecs=10
priority=999
